# templates/default/tram-upstart.conf.erb
description "tram webapp"
author "Nick Gerakines, nick@gerakines.net"

start on (net-device-up and local-filesystems and runlevel [2345])
stop on (runlevel [016])

chdir /home/tram/
exec su tram -c '/home/tram/tram'
